# config/systemd/user/getmail.service
[Unit]
Description=Run getmail to get new email

[Service]
ExecStart=/usr/bin/getmail --verbose

[Install]
WantedBy=default.target
